# deploy/s6/migrate/up
/command/with-contenv /app/bin/migrate.sh
